# webhook/webhook.service
[Unit]
Description=Universe Webhook Receiver — GitHub → Telegram
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/theuniverse
EnvironmentFile=/opt/oraculo/.env
ExecStart=/opt/webhook/venv/bin/uvicorn webhook.receiver:app --host 0.0.0.0 --port 9120
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook

[Install]
WantedBy=multi-user.target
