# obi-wan/obi-wan.service
[Unit]
Description=Obi-Wan — conselheiro conversacional do TheUniverse
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/theuniverse
EnvironmentFile=/opt/obi-wan/.env
ExecStart=/opt/obi-wan/venv/bin/python /opt/theuniverse/obi-wan/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obi-wan

[Install]
WantedBy=multi-user.target
